# deploy/vfs-check.timer
[Unit]
Description=Run the VFS appointment check about every 20 minutes

[Timer]
# First run shortly after boot, then ~every 20 min measured from the end of
# the previous run. RandomizedDelaySec adds OS-level jitter on top of check.py's
# own 0-240s in-process jitter, so polls are well de-synchronised from the
# clock (anti-detection) and never overlap.
OnBootSec=3min
OnUnitActiveSec=30min
RandomizedDelaySec=180
Unit=vfs-check.service

[Install]
WantedBy=timers.target
